# parser/crontab
# Получение нового контента каждый день
* * * * * /usr/local/bin/node /app/index.js
